# config_files/cron/crontab
@hourly /home/roman/send_ftp.py

/6 * * * * /home/roman/incomming_script.sh > /var/log/cron.log
