# config_files/supervisor/wvdial.conf
[program:wvdial]
command=wvdial
autostart=true
autorestart=true
stderr_logfile=/var/log/wvdial_error.log
stdout_logfile=/var/log/wvdial_out.log
stopsignal=KILL
numprocs=1 